# services/applications/mopro/supervisord.mopro_sh.conf
[program:mopro]
priority=99999
command=bash /opt/qcrbox/run_once.sh
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
