# services/applications/olex2_linux/supervisord.olex2_headless.conf
[program:olex2-headless]
priority=1
command=/opt/olex2/startc
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
